# docker/etc.slurm/services.d/munge/run
#!/usr/bin/execlineb -P

s6-setuidgid munge
/usr/sbin/munged --num-threads=1 -F
